# auto-restart-openvpn.service
[Unit]
Description=Startup script auto-restart-openvpn

[Service]
ExecStart=/usr/bin/auto-restart-openvpn
